# deploy/systemd/salini-api.service
[Unit]
Description=Salini AMS API (.NET Core)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/salini-ams/backend/salini.api.API
ExecStart=/usr/bin/dotnet /var/www/salini-ams/backend/salini.api.API/salini.api.API.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=salini-api
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/salini-ams/backend/salini.api.API/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
